# systemd/public-terminal-rebalance.timer
[Unit]
Description=Run S&P 500 rebalance daily at 12:00 ET

[Timer]
# Fire at noon Eastern Time every weekday (market is open Mon–Fri)
OnCalendar=Mon..Fri *-*-* 12:00:00
TimeZone=America/New_York
# If the system was off at noon, run as soon as it comes back up
Persistent=true

[Install]
WantedBy=timers.target
